# omochao.service.example
[Unit]
Description=Omochao Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/youruser/projects/omochao
ExecStart=/home/youruser/.venv/omochao/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
